# modules/headphones/templates/initd.erb
#!/sbin/runscript
#
# DO NOT EDIT
#
# This file was generated by Puppet
# See the 'headphones' Puppet module for details                                                                                                                                                                                              
#

HEADPHONES_PID_DIR="/var/run/headphones"

command="/usr/bin/python2.7"
command_args="<%= @home %>/Headphones.py --nolaunch --daemon --pidfile ${HEADPHONES_PID_DIR}/headphones.pid"
pidfile="${HEADPHONES_PID_DIR}/headphones.pid"
start_stop_daemon_args="--user headphones"
name="Headphones"
description="Automatic Music Downloader"

depend() {
    need localmount net
    use sabnzbd transmission
}

start_pre() {
    checkpath --directory --owner headphones:headphones "${HEADPHONES_PID_DIR}"
}
