# modules/couchpotato/templates/initd.erb
#!/sbin/runscript
#
# DO NOT EDIT
#
# This file was generated by Puppet
# See the 'couchpotato' Puppet module for details                                                                                                                                                                                              
#

COUCHPOTATO_PID_DIR="/var/run/couchpotato"

command="/usr/bin/python2.7"
command_args="<%= @home %>/CouchPotato.py --daemon --pid_file ${COUCHPOTATO_PID_DIR}/couchpotato.pid"
pidfile="${COUCHPOTATO_PID_DIR}/couchpotato.pid"
start_stop_daemon_args="--user couchpotato"
name="CouchPotato"
description="Automatic Movie Downloading via NZBs & Torrents"

depend() {
    need localmount net
    use sabnzbd
}

start_pre() {
    checkpath --directory --owner couchpotato:couchpotato "${COUCHPOTATO_PID_DIR}"
}
