# tmp/nest/files/bluetooth/bluetooth-sleep.service
[Unit]
Description=Bluetooth sleep hook
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=-/usr/bin/systemctl stop bluetooth
ExecStop=-/usr/bin/systemctl start bluetooth

[Install]
WantedBy=sleep.target
