# modules/crashplan/files/crashplan.init
#!/sbin/runscript
#
# DO NOT EDIT
#
# This file was generated by Puppet
# See the 'crashplan' Puppet module for details
#

depend() {
    need net
    need localmount
}

start() {
    ebegin "Starting CrashPlan"
    /usr/local/crashplan/bin/CrashPlanEngine start > /dev/null
    eend $?
}

stop() {
    ebegin "Stopping CrashPlan"
    /usr/local/crashplan/bin/CrashPlanEngine stop > /dev/null
    eend $?
}
