# files/streamux/mnt-external.automount
[Unit]
Before=local-fs.target

[Automount]
Where=/mnt/external

[Install]
WantedBy=local-fs.target
