# modules/nzbdrone/templates/initd.erb
#!/sbin/runscript
#
# DO NOT EDIT
#
# This file was generated by Puppet
# See the 'nzbdrone' Puppet module for details                                                                                                                                                                                              
#

command="/usr/bin/mono"
command_args="<%= @home %>/NzbDrone.exe"
pidfile="/var/run/nzbdrone.pid"
start_stop_daemon_args="--user nzbdrone"
command_background="true"
name="NzbDrone"
description="Smart PVR for newsgroup users."

depend() {
    need localmount net
    use sabnzbd
}
