# files/ydotool/ydotoold.socket
[Socket]
ListenStream=/run/ydotool.sock
SocketGroup=users
SocketMode=0660
