# files/ydotool/ydotoold.socket
[Socket]
ListenStream=/tmp/.ydotool_socket
SocketGroup=users
SocketMode=0660
